# relay/gertec-relay.service
[Unit]
Description=GERTEC Firebase Notification Relay
After=network.target

[Service]
Type=simple
User=gertec
WorkingDirectory=/opt/gertec-relay
ExecStart=/usr/bin/node server.js
EnvironmentFile=/opt/gertec-relay/.env
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
